# rasa/crontab
0 7 * * 1-5 /notify.sh >> /home/out_cron.log 2>&1
* * * * *  /notify_gmail.sh >> /home/out_cron.log 2>&1
