# rasa/crontab
* * * * * curl webcrawler:5000/cardapio/Monday > /home/out_cron.json
